# cron/.config/cron/crontab
# run `crontab $HOME/.config/cron/crontab` to apply changes
*/5 * * * * $HOME/.local/bin/battery-alert
